# config/systemd/labctl-web.service
[Unit]
Description=Lab Controller Web Dashboard
Documentation=https://github.com/johnjezl/Embedded-Lab-Control
After=network.target ser2net.service

[Service]
Type=simple
User=labctl
Group=labctl
ExecStart=/opt/labctl/venv/bin/labctl -c /etc/labctl/config.yaml web --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=10

# For HTTPS, configure cert_file and key_file in the config.yaml
# under the 'web:' section. See config/labctl.yaml.example.
# To bind to port 443, uncomment AmbientCapabilities below.
# AmbientCapabilities=CAP_NET_BIND_SERVICE

# Environment and paths
Environment="HOME=/var/lib/labctl"
Environment="FLASK_ENV=production"
WorkingDirectory=/var/lib/labctl

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=labctl-web

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/labctl
ReadOnlyPaths=/etc/labctl/ssl

[Install]
WantedBy=multi-user.target
