# config/systemd/labctl-web.service
[Unit]
Description=Lab Controller Web Dashboard
Documentation=https://github.com/johnjezl/Embedded-Lab-Control
After=network.target ser2net.service

[Service]
Type=simple
User=labctl
Group=labctl
ExecStart=/opt/labctl/venv/bin/labctl web --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=10

# Environment and paths
Environment="HOME=/var/lib/labctl"
Environment="FLASK_ENV=production"
WorkingDirectory=/var/lib/labctl

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=labctl-web

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/labctl

[Install]
WantedBy=multi-user.target
